# web-frontend/gunicorn-gpioneer.conf
[program:gunicorn]
command=sudo gunicorn app:app -b 0.0.0.0:80
autostart=true
autorestart=false
stderr_logfile=/var/log/long.err.log
stdout_logfile=/var/log/long.out.log